# deploy/wallac-designer.service
[Unit]
Description=Wallac Victor2 Protocol Designer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=antonio
WorkingDirectory=/home/antonio/repos/wallac-victor2-api
EnvironmentFile=/etc/wallac-bridge/bridge.env
ExecStart=/home/antonio/repos/wallac-victor2-api/.venv/bin/python3 -m uvicorn bridge.designer_app:create_designer_app --factory --host 0.0.0.0 --port 8422

Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
