# deploy/wallac-bridge.service
[Unit]
Description=Wallac Victor2 eLabFTW Bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=antonio
WorkingDirectory=/home/antonio/repos/wallac-victor2-api
EnvironmentFile=/etc/wallac-bridge/bridge.env
ExecStart=/home/antonio/repos/wallac-victor2-api/.venv/bin/python3 -m uvicorn bridge.bridge_app:create_bridge_app --factory --host 0.0.0.0 --port 8423

Restart=on-failure
RestartSec=10

# Security hardening. The bridge only writes to the spool dir; everything
# else (code, venv, env file) is read-only.
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/lib/wallac-bridge
ProtectHome=read-only

[Install]
WantedBy=multi-user.target
